# cron/cron.daily
#!/bin/bash

# cron/cron.daily - Script executed once per day by LoxBerry's cron system.
#
# Installed to $LBHOMEDIR/system/cron/cron.daily/<NAME> during plugin installation.
# Runs as user "loxberry".

# Add your daily task here.
