# cron/cron.monthly
#!/bin/bash

# cron/cron.monthly - Script executed once per month by LoxBerry's cron system.
#
# Installed to $LBHOMEDIR/system/cron/cron.monthly/<NAME> during plugin installation.
# Runs as user "loxberry".

# Add your monthly task here.
